# dist/etc/systemd/system/paperframe.service
[Unit]
Description=Paperframe Service
After=raspi-config.service
Wants=network-online.target

[Service]
Type=exec
ExecStart=/bin/paperframe service
RemainAfterExit=no
SuccessExitStatus=0
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
